# .docker/staging/fez.cron
# Example of job definition:
# .---------------- minute (0 - 59)
# |  .------------- hour (0 - 23)
# |  |  .---------- day of month (1 - 31)
# |  |  |  .------- month (1 - 12) OR jan,feb,mar,apr ...
# |  |  |  |  .---- day of week (0 - 6) (Sunday=0 or 7) OR sun,mon,tue,wed,thu,fri,sat
# |  |  |  |  |
# *  *  *  *  * user-name  command to be executed

## Fez staging tasks ##

# Update the GeoIP and GeoLiteCity data files in S3 every month on the 3rd at 3:33am
33  3   3 * *   root    /usr/bin/wget -nd -q -O /tmp/GeoIP.dat.gz http://www.maxmind.com/download/geoip/database/GeoIP.dat.gz && aws s3 cp /tmp/GeoIP.dat.gz s3://uql-fez-staging-cache/GeoIP.dat.gz > /dev/null 2>&1
33  3   3 * *   root    /usr/bin/wget -nd -q -O /tmp/GeoLiteCity.dat.gz http://www.maxmind.com/download/geoip/database/GeoLiteCity.dat.gz && aws s3 cp /tmp/GeoIP.dat.gz s3://uql-fez-staging-cache/GeoLiteCity.dat.gz > /dev/null 2>&1

# Run the integrity fixing scripts each morning
15  8   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_run_integrity_checks&class=BackgroundProcess_Run_Integrity_Checks&input[runType]=both&token=WEBCRON_TOKEN" > /dev/null 2>&1
45  8   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_run_integrity_checks&class=BackgroundProcess_Run_Integrity_Checks&input[runType]=fix&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Refresh Solr with any pids changed in the last 5 months (only needed in staging for espace nightly rebuilds)
00 7 * * *  root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_reindex_recent_pids_into_solr&class=BackgroundProcess_Reindex_Recent_Pids_Into_Solr&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Update statistics summary tables each day
30 13   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_statistics_summary_tables&class=BackgroundProcess_Update_Statistics_Summary_Tables&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Cache table rebuild
30 13          * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_cache_rebuild&class=BackgroundProcess_Cache_Rebuild&token=WEBCRON_TOKEN" > /dev/null 2>&1
00 1-6,18-23   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_cache_rebuild&class=BackgroundProcess_Cache_Rebuild&input[runRebuild][]=CloudTags&input[runRebuild][]=RecentDownloads&input[recentDownloadsLimit]=5&token=WEBCRON_TOKEN" > /dev/null 2>&1
00         *   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_cache_rebuild&class=BackgroundProcess_Cache_Rebuild&input[runRebuild][]=RecentItems&input[recentItemsLimit]=5&token=WEBCRON_TOKEN" > /dev/null 2>&1

# TR Web Services Premium Wok Queue processor
#0 1     * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_process_wok_queue&class=BackgroundProcess_Process_Wok_Queue&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Get citation count from WoS and Scopus
#0 0     * * 0   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_citation_counts&class=BackgroundProcess_Update_Citation_Counts&input[runType]=wsp&token=WEBCRON_TOKEN" > /dev/null 2>&1
#7 2     * * 0   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_citation_counts&class=BackgroundProcess_Update_Citation_Counts&input[runType]=scopus&token=WEBCRON_TOKEN" > /dev/null 2>&1
#0 0     * * 0   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_citation_counts&class=BackgroundProcess_Update_Citation_Counts&input[runType]=scopus-empty&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Download new UQ pubs from WoS / Scopus / ResearcherID
#0 0     * * 0   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_download_uq_pubs&class=BackgroundProcess_Download_Uq_Pubs&input[runType]=wos&token=WEBCRON_TOKEN" > /dev/null 2>&1
#0 0     * * 2   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_download_uq_pubs&class=BackgroundProcess_Download_Uq_Pubs&input[runType]=scopus&token=WEBCRON_TOKEN" > /dev/null 2>&1
#4 18    * * 5   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_download_uq_pubs&class=BackgroundProcess_Download_Uq_Pubs&input[runType]=researcherid&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Run Links AMR check
#0 0     * * 3   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_links_amr_check&class=BackgroundProcess_Links_Amr_Check&token=WEBCRON_TOKEN" > /dev/null 2>&1

# ERA Ranked Journal matching
#0 19    * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_match_ranked_journals&class=BackgroundProcess_Match_Ranked_Journals&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Link Status Report
#0 0     1 * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_check_links&class=BackgroundProcess_Check_Links&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Delete downloaded XML content for ResearcherID Profile & Publications that are older than 2 weeks
#0 2     * * 1   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_clean_researcherid_xml&class=BackgroundProcess_Clean_Researcherid_Xml&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Re-download Sherpa Romeo data in case it has updated
#0 15    * * 6   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_sherpa_romeo_data&class=BackgroundProcess_Update_Sherpa_Romeo_Data&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Download Altmetric scores using the Altmetric API
#0 1     * * 1   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_altmetric_info&class=BackgroundProcess_Update_Altmetric_Info&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Open up open access documents past embargo
#15 18   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_embargo_period_complete_check&class=BackgroundProcess_Embargo_Period_Complete_Check&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Re-download Ulrichs data in case it has updated
#0 15    * * 6   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_ulrichs&class=BackgroundProcess_Update_Ulrichs&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Update pids with Ulrichs data
#0 15    * * 1   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_oa_pids_with_ulrichs&class=BackgroundProcess_Update_Oa_Pids_With_Ulrichs&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Look for missing Scopus IDs given DOI
#0 5     * * 6   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_insert_scopus_id_using_doi_search&class=BackgroundProcess_Insert_Scopus_Id_Using_Doi_Search&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Find auto refereed matches
#30 6    * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_set_refereed_details&class=BackgroundProcess_Set_Refereed_Details&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Find all EIDs that have expired
#30 18   * * 1   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_find_expired_eids&class=BackgroundProcess_Find_Expired_Eids&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Check for datastream policy conflicts. It could cause security issues so will need an alert ASAP.
#0 0-23/2 * * *  root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_alert_datastream_policy_conflicts&class=BackgroundProcess_Alert_Datastream_Policy_Conflicts&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Update all Publons reviews
#25 7    * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_update_all_publons_reviews&class=BackgroundProcess_Update_All_Publons_Reviews&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Update Sitemap
#30 18   * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_generate_sitemap&class=BackgroundProcess_Generate_Sitemap&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Email users who have subscribed searches
#30 4    * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_email_subscribed_search_users&class=BackgroundProcess_Email_Subscribed_Search_Users&token=WEBCRON_TOKEN" > /dev/null 2>&1

# eSpace / Eventum job sync
#*/15 *  * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_sync_eventum_jobs&class=BackgroundProcess_Sync_Eventum_Jobs&token=WEBCRON_TOKEN" > /dev/null 2>&1

# ResearcherID Batch Download - Checks the status of jobs created using the ResearcherID Batch Download service
#*/5 0-23 * * *  root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_check_researcherid_download_status&class=BackgroundProcess_Check_Researcherid_Download_Status&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Process mail queue
#*/2 *    * * *  root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_process_mail_queue&class=BackgroundProcess_Process_Mail_Queue&token=WEBCRON_TOKEN" > /dev/null 2>&1

# ResearcherID Batch Upload - Process reports received from the service
#*/5 0-23 * * *  root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_process_researcherid_upload_emails&class=BackgroundProcess_Process_Researcherid_Upload_Emails&token=WEBCRON_TOKEN" > /dev/null 2>&1

# Import new RHD students from sinet_load database
#0 1     * * *   root    /usr/bin/wget -O - "http://fez-staging.library.uq.edu.au/api/cron_register_bgp.php?file=bgp_import_new_rhd_students&class=BackgroundProcess_Import_New_Rhd_Students&token=WEBCRON_TOKEN" > /dev/null 2>&1

## /Fez staging tasks ##

